# deploy/systemd/obs-live-studio-stream-relay.service
[Unit]
Description=OBS Live Studio YouTube/Twitch Stream Relay
Wants=network-online.target
After=network-online.target
Before=obs-live-studio-desktop-agent.service

[Service]
Type=simple
WorkingDirectory=%h/obs-live-studio
EnvironmentFile=%h/obs-live-studio/.env
ExecCondition=/bin/sh -c '[ "${MULTISTREAM_ENABLED:-false}" = "true" ]'
ExecStartPre=/usr/bin/env node --env-file=%h/obs-live-studio/.env %h/obs-live-studio/scripts/configure-stream-relay.mjs
ExecStart=%h/obs-live-studio/scripts/run-stream-relay.sh
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/obs-live-studio/var

[Install]
WantedBy=obs-live-studio.target
