# deploy/systemd/obs-live-studio-worker.service
[Unit]
Description=OBS Live Studio News Worker
After=network-online.target obs-live-studio-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/obs-live-studio
EnvironmentFile=-%h/obs-live-studio/.env
Environment=NODE_ENV=production
ExecStart=/bin/bash -lc 'source "$HOME/.nvm/nvm.sh"; exec npm run start -w @ans/worker'
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=default.target
