# deploy/systemd/obs-live-studio-broadcast-runner.service
[Unit]
Description=OBS Live Studio Broadcast Runner
After=obs-live-studio-api.service obs-live-studio-desktop-agent.service
Wants=obs-live-studio-api.service obs-live-studio-desktop-agent.service

[Service]
Type=simple
WorkingDirectory=%h/obs-live-studio
EnvironmentFile=-%h/obs-live-studio/.env
Environment=NODE_ENV=production
ExecStart=/bin/bash -lc 'source "$HOME/.nvm/nvm.sh"; exec npm run start -w @ans/broadcast-runner'
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=default.target
